# patch_watcher/patch_watcher_daily.timer
[Unit]
Description=Lustre Patch Watcher Daily Confidence Report Timer

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
