# patch_watcher/patch_watcher.timer
[Unit]
Description=Lustre Patch Watcher Timer

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
