# patch_shepherd/patch_shepherd_daily.timer
[Unit]
Description=Lustre Patch Shepherd Daily Confidence Report Timer

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
