# patch_shepherd/patch_shepherd_daily.service
[Unit]
Description=Lustre Patch Shepherd Daily Confidence Report
After=network-online.target

[Service]
Type=oneshot
User=patchshepherd
ExecStart=/shared/llm_code_and_review_tools/patch_shepherd/daily_confidence.sh
TimeoutStartSec=120
Environment=HOME=/home/patchshepherd
Environment=PATH=/home/patchshepherd/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/shared/llm_code_and_review_tools/patch_shepherd
